# etc/cron.d/parse_ppp_for_fail2ban
MAILTO=""
*	*	*	*	*	root	/usr/bin/parse_pptp_log.py /var/log/ppp > /var/log/ppp.login
